# src/log/logrotate.cron
00 0 * * * /usr/bin/logrotate -s /tmp/kylog-rotate-tmp /etc/kysdk/kysdk-base/kylog-rotate-default

